# testing/k2http-supervisor.conf
[program:k2http]
command=/app/k2http --config /app/config.yml --debug
autostart=true
autorestart=true
